# systemd/check_pritunl.timer
[Unit]
Description=Run Pritunl check every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
