# systemd/check_pritunl.service
[Unit]
Description=Check Pritunl status


[Service]
Type=oneshot
ExecStart=/opt/automation/scripts/check_pritunl.sh
